# deb/build/pzomboidService.0.0-1/etc/systemd/system/pzomboid.service
[Unit]
Description=Handles Project Zomboid Server Application
Wants=syslog.target network.target network-online.target

[Service]
Type=simple
User=pzuser
Group=pzuser
WorkingDirectory=/home/pzuser
ExecStart=/bin/bash /opt/pzserver/pzserver-start /opt/pzserver devServer01
ExecStop=/bin/bash /opt/pzserver/pzserver-stop /home/pzuser
Sockets=pzomboid.socket
StandardInput=socket
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target